# http-s/supervisor/chat-ws.conf
[program:chat-ws]
command=/usr/bin/python3 /opt/chat/server.py
directory=/opt/chat
autostart=true
autorestart=true
stdout_logfile=/var/log/chat-ws.log
stderr_logfile=/var/log/chat-ws.err
